# webserver_service/flask.service
[Unit]
Description=Start Flask server

[Service]
ExecStart=/usr/bin/flask

[Install]
WantedBy=multi-user.target
