# webserver_service/startflask.service
[Unit]
Description=Start Flask server

[Service]
ExecStart=/usr/bin/startflask

[Install]
WantedBy=multi-user.target
